# deb/dokku-alt/etc/init/dokku-daemon.conf
description "Dokku daemon service"

start on (started docker)

script
  sleep 10 # give us some time
  /usr/local/bin/dokku-daemon
end script
