# deb/dokku-alt/etc/init/dokku-redeploy.conf
description "Dokku app redeploy service"

start on (started docker)

script
  sleep 10 # give docker some time
  /usr/local/bin/dokku deploy:all
end script
